# profile.d/aliases.ps1
Set-Alias -Name npm -Value pnpm
Set-Alias -Name rmf -Value Remove-Fast
Set-Alias -Name which -Value Get-Command
Set-Alias -Name condaaa -Value Start-Conda
Set-Alias -Name lg -Value lazygit
Set-Alias -Name rgaf -Value Search-RgaFzf
Set-Alias -Name psf -Value Format-PowerShellScript
Set-Alias -Name ib -Value Invoke-Build
