# profile.d/functions.ps1
Get-ChildItem -Path ".\functions.d" -Filter "*.ps1"
| ForEach-Object {.$_ }
