# profile.d/modules.ps1
Import-Module "${SCOOP_HOME}\modules\scoop-completion"
Import-Module "${SCOOP_HOME}\modules\lazy-posh-git"
Import-Module -Name CompletionPredictor
