# profile.d/variables.ps1
$SCOOP_HOME = "D:\Scoop"
$env:ASTRO_CONFIG = "${env:LOCALAPPDATA}\nvim\"
$env:ASTRO_USER_CONFIG = "${env:LOCALAPPDATA}\nvim\lua\user"
